# zaptec-solarcharge.service
[Unit]
Description=Zaptec Solarcharge - Zonnepanelen EV Laadregeling
Documentation=https://github.com/materharmsel/zaptec-solarcharge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/zaptec-solarcharge
ExecStart=/home/pi/zaptec-solarcharge/venv/bin/python main.py
Restart=on-failure
RestartSec=30

# Stdout en stderr gaan naar het systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zaptec-solarcharge

[Install]
WantedBy=multi-user.target
